# config/systemd/vllm-paul-voice.service.j2
# Managed by book-pipeline. Regenerate via `book-pipeline vllm-bootstrap`.
# Serves the pinned voice-FT LoRA adapter on {{ port }}. Do NOT edit by hand —
# edits are overwritten on the next bootstrap run.

[Unit]
Description=vLLM serving Paul voice-FT LoRA adapter ({{ ft_run_id }}) on port {{ port }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
Environment=VLLM_USE_V1=0
Environment=HF_HUB_DISABLE_XET=1
# EnvironmentFile= dropped per Forge inspection 2026-04-24 — the only file
# available is /home/admin/finetuning/cu130_env.sh which is shell-script
# format (export FOO=BAR), unparseable by systemd's KEY=VALUE expectation.
# Inline Environment= directives above cover the critical vLLM env vars.
ExecStartPre=/bin/sh -c 'nvidia-smi > /dev/null 2>&1 || exit 1'
ExecStart={{ venv_python }} -m vllm.entrypoints.openai.api_server \
    --model {{ base_model }} \
    --enable-lora \
    --lora-modules paul-voice={{ adapter_path }} \
    --served-model-name paul-voice \
    --host 127.0.0.1 \
    --port {{ port }} \
    --dtype {{ dtype }} \
    --max-model-len {{ max_model_len }} \
    --max-num-seqs 4 \
    --tensor-parallel-size {{ tensor_parallel_size }} \
    --gpu-memory-utilization {{ gpu_memory_utilization }} \
{% if quantization and quantization != 'none' %}    --quantization {{ quantization }} \
{% endif %}    --enforce-eager \
    --trust-remote-code \
    --default-chat-template-kwargs '{"enable_thinking": false}'
Restart=on-failure
RestartSec=15
TimeoutStartSec=180

[Install]
WantedBy=default.target
